# deploy/systemd/staysense-import.service
[Unit]
Description=StaySense Data Import Job
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=staysense
Group=staysense
WorkingDirectory=/opt/staysense/backend
UMask=0002
ExecStart=/usr/bin/python3 /opt/staysense/backend/run_import_jobs.py --config /opt/staysense/docs/open_data_sources_nrw_live.json --prune-legacy
